# src/create_flask/templates/deploy/supervisor/{{ project_name }}.conf
; supervisor 进程配置（部署前请按实际环境修改 directory / user / 路径）
[program:{{ project_name }}]
directory=/srv/{{ project_name }}
command=/srv/{{ project_name }}/.venv/bin/gunicorn -c gunicorn.conf.py wsgi:app
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/{{ project_name }}/logs/gunicorn.log
environment=ENV="production"
{% if use_celery %}
[program:{{ project_name }}-worker]
directory=/srv/{{ project_name }}
command=/srv/{{ project_name }}/.venv/bin/celery -A wsgi:celery_app worker -l info
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/{{ project_name }}/logs/celery.log
environment=ENV="production"
{% endif %}
